# cpp/deploy/rest_server/rest_server.conf
[program:rest_server]
command=/prod/run.sh
stdout_logfile=/var/log/supervisor.log.stdout
stderr_logfile=/var/log/supervisor.log.stderr
autostart=true
autorestart=true
